# _provision/roles/packages/templates/unbound/net.unbound.unbound.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-/Apple/DTD PLIST 1.0/EN" "http:/www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>net.unbound.unbound</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/sbin/unbound</string>
      <string>-d</string>
      <string>-c</string>
      <string>/usr/local/svcs/unbound/etc/unbound.conf</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/usr/local/svcs/unbound/var/unbound.log</string>
    <key>StandardOutPath</key>
    <string>/usr/local/svcs/unbound/var/unbound.log</string>
  </dict>
</plist>
