# sv/runsvdir/run
#!/bin/sh

_user=${PWD##*-}
_groups=$(id -Gn "$_user" | tr ' ' ':')

if [ ! -d "/run/runit.$_user" ]; then
    install -d -o"$_user" -g"$_user" "/run/runit.$_user"
fi

USER="$_user"; export USER
HOME="/home/$_user"; export HOME

exec 2>&1
exec chpst -u "$_user:$_groups" runsvdir "$HOME/.local/var/service" 'log: ...........................................................................................................................................................................................................................................................................................................................................................................................................'
