# provision/roles/packages/templates/home/Library/LaunchAgents/th.in.grid.env.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>th.in.grid.env</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/fish</string>
    <string>-l</string>
    <string>-c</string>
    <string>
      launchctl setenv SSH_AUTH_SOCK $SSH_AUTH_SOCK
      launchctl setenv SSH_AGENT_PID $SSH_AGENT_PID
      launchctl setenv GPG_AGENT_INFO $GPG_AGENT_INFO
    </string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
