# sv/podman/run
#!/usr/bin/env -S execlineb -P

emptyenv
export PATH /usr/local/bin:/usr/bin:/bin

# Setting up
#

backtick -n user { id -un }
backtick -n uid { id -u }

multisubstitute {
    importas -u user user
    importas -u uid uid
}

backtick -n homedir { homeof $user }

multisubstitute {
    importas -u homedir homedir
    define xdg-runtime-dir /run/user/${uid}
}

# Only proceed if XDG_RUNTIME_DIR already created by elogind. Without this
# check podman will restart when we first login which is not what we want.
#

if { test -d ${xdg-runtime-dir} }

export HOME ${homedir}
export XDG_RUNTIME_DIR ${xdg-runtime-dir}

fdmove -c 2 1
/usr/bin/podman system service -t 0 --log-level=debug unix://${xdg-runtime-dir}/podman/podman.sock
